# db/backup.cron
0 17 * * 1,2,3,4,5 /usr/local/bin/docker exec -it m0ney-db /usr/bin/mysqldump --databases stocks | gz > /home/m0ney/db/data/$(date +%Y-%m-%d).sql.gz && /usr/local/bin/docker exec -it m0ney-db mysqldump --databases stocks > /home/m0ney/db/data/latest.sql
